# crontab
0 * * * * /bin/bash /usr/src/app/wrapper.sh >> /usr/src/app/app/logs/cron.log 2>&1
